# publisher/blog.cron.d

# set config options in /var/www/publisher/config

* * * * * www-data /bin/bash /var/www/publisher/publisher.sh
